# daemon/kimi-github-agent.service
[Unit]
Description=Kimi GitHub Agent
After=network.target

[Service]
Type=simple
User=abds
WorkingDirectory=/home/abds/clawd/Kimi-github-agent
Environment=NODE_ENV=production
Environment=PORT=3001
ExecStart=/usr/bin/node /home/abds/clawd/Kimi-github-agent/src/index.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kimi-github-agent

[Install]
WantedBy=multi-user.target
